# dist/opnsense-update-notify.service
[Unit]
Description=An update notification script for OPNsense.
Wants=opnsense-update-notify.timer
After=network-online.target
Requires=network-online.target

[Service]
WorkingDirectory=/etc/opnsense-update-notify
ExecStart=/usr/bin/python3 /usr/bin/opnsense-update-notify
Restart=always

[Install]
WantedBy=multi-user.target
